# templates/fetch_secrets.service
[Unit]
Description=Fetch secrets from Key Vault into /run (managed identity)
Before=nginx.service ask-app.service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/fetch-secrets.sh

[Install]
WantedBy=multi-user.target
