# templates/mcp_db.service.j2
[Unit]
Description=MCP Infrastructure Server (db-vm)
After=mariadb.service
Documentation=https://modelcontextprotocol.io

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mcp-db
Environment=MCP_DB_NAME={{ db_name }}
ExecStart=/opt/mcp-db/bin/python /opt/mcp-db/mcp_db.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
