# templates/sentinel.service
[Unit]
Description=Sentinel health check
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sentinel.sh
StandardOutput=journal
StandardError=journal
