# minecraft/mc_guide.service
[Unit]
Description=FTB NeoTech GuideBot API
After=network.target

[Service]
Type=simple
User=jrick
WorkingDirectory=/opt/mc-guide
EnvironmentFile=/opt/mc-guide/.env
ExecStart=/opt/mc-guide/venv/bin/uvicorn mc_guide:app --host 127.0.0.1 --port 8765
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
